# configure/crontab
# m h  dom mon dow   command

* 1 * * * /exfe/exfebus/bin/clear.sh > /var/log/exfebus/clear.log
